# deploy/lolbench.timer.j2
[Unit]
Description=lolbench agent
After=syslog.target network.target

[Timer]
OnBootSec=300
Persistent=true

[Install]
WantedBy=multi-user.target
