# deploy/lolbench.service.j2
[Unit]
Description=lolbench agent

[Service]
Type=oneshot
SendSIGKILL=no
ExecStart={{ deploy_dir }}/process.sh
