# deploy/lolbench.timer.j2
[Unit]
Description=lolbench agent
After=network.target
Requires=lolbench.service

[Timer]
OnUnitActiveSec=30
Persistent=true

[Install]
WantedBy=timers.target
